# setup/conf/site-cron.crontab
# Run cron every five minutes
*/5 * * * * cd $DRUPAL_SITE_DIR && drush cron 
